# seaboat-gunicorn.service
[Unit]
Description=SeaBoat Gunicorn
After=network.target mysqld.service

[Service]
User=www
Group=www
WorkingDirectory=/www/wwwroot/seabot/back-end/seadrone
# 改为你实际的虚拟环境 PATH
Environment="PATH=/www/wwwroot/seabot/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=seadrone.settings_prod"
Environment="PROJECT_ROOT=/www/wwwroot/seabot"
Environment="SHIP_ONLINE_GRACE_SECONDS=30"
EnvironmentFile=/www/wwwroot/seabot/.env
# ExecStart 中的 gunicorn 路径改为项目内虚拟环境的可执行文件
ExecStart=/www/wwwroot/seabot/venv/bin/gunicorn \
  --workers 4 \
  --bind 127.0.0.1:8000 \
  --access-logfile /www/wwwroot/seabot/logs/gunicorn_access.log \
  --error-logfile /www/wwwroot/seabot/logs/gunicorn_error.log \
  seadrone.wsgi:application
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
