# seabot-gateway.service
[Unit]
Description=SeaBoat TCP Gateway (9001-9004)
After=network.target

[Service]
Type=simple
User=www
Group=www
WorkingDirectory=/www/wwwroot/seabot/back-end/seadrone
Environment="PATH=/www/wwwroot/seabot/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=seadrone.settings_prod"
EnvironmentFile=-/www/wwwroot/seabot/.env
ExecStart=/www/wwwroot/seabot/venv/bin/python manage.py run_ship_gateway --host 0.0.0.0 --port-start 9001 --ship-count 2 --status-interval 30
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target